# build/debian/postinst
#!/bin/sh

#
# npemwin user/group, directories owned by user, and directories owned by root
#
user=noaaport
home=/var/noaaport
#
userdirs="/var/run/npemwin
/var/log/npemwin
/var/log/npemwin/www
/var/npemwin
/var/npemwin/stats
/var/npemwin/spool
/var/npemwin/spool/data
/var/npemwin/spool/lock
/var/npemwin/spool/tmp
/var/npemwin/spool/qfiles
/var/npemwin/files
/var/npemwin/inv
/var/npemwin/rss
/var/npemwin/data
/var/npemwin/data/digatmos
/var/npemwin/data/latest
/var/npemwin/data/metar
/var/npemwin/data/msg
/var/npemwin/data/gempak
/var/npemwin/data/track
/var/npemwin/archive
/var/npemwin/archive/data
/var/npemwin/archive/data/digatmos
/var/npemwin/archive/data/metar
/var/npemwin/archive/data/spool
/var/npemwin/archive/inv
/var/npemwin/archive/inv/digatmos
/var/npemwin/www
/var/npemwin/tmp"
#
rootdirs="/usr/local/etc/npemwin/defaults
/usr/local/etc/npemwin/site
/usr/local/etc/npemwin/rc.d
/usr/local/etc/npemwin/rc.d/std
/usr/local/etc/npemwin/rc.d/metar
/usr/local/etc/npemwin/rc.d/da
/usr/local/etc/npemwin/rc.d/gp
/usr/local/etc/npemwin/rc.d/ldm
/usr/local/etc/npemwin/rc.d/nntp
/usr/local/etc/npemwin/rc.d/rss
/usr/local/etc/npemwin/rc.d/npemwin
/usr/local/etc/npemwin/rc.d/pan
/usr/local/etc/npemwin/rc.d/msg
/usr/local/libexec/npemwin/site
/usr/local/libexec/npemwin/tclhttpd/site"
#
# target directories for optional files
#
docdir=/usr/local/share/doc/npemwin
confdir=/usr/local/etc/npemwin
rcdir=/etc/init.d
defaultdir=/etc/default

#
# noaaport user
#
useradd -d $home $user > /dev/null 2>&1
status=$?
if [ $status -ne 9 -a $status -ne 0 ]
then
    echo "Error creating user $user."
    exit 1
fi

for d in $userdirs
do
  [ -e $d ] || install -d -m 755 -o $user -g $user $d
done

for d in $rootdirs
do
  [ -e $d ] || install -d -m 755 $d
done

#
# install optional files from doc directory
#
install -m 0755 $docdir/npemwind.sh-sample $rcdir/npemwind
[ -f $defaultdir/npemwind ] || \
    install -m 0644 $docdir/npemwind.sh.conf-debian $defaultdir/npemwind

if [ -n "`which chkconfig`" ]
then
	chkconfig --add npemwind
elif [ -n "`which update-rc.d`" ]
then
	update-rc.d npemwind defaults 99
fi

#
# configure default
#
cd $confdir/dist
for file in *.conf-sample *.rc-sample *.def-sample
do
  name=${file%%"-sample"}
  cp $file ../$name
done
cp README ..

for file in *-defaults
do
  name=${file%%"-defaults"}
  cp $file ../defaults/$name
done

for file in *-templates.tgz
do
  tar -C ../defaults -xzf $file
done

#
# www and tclhttpd
#
cp /usr/local/share/doc/npemwin/www.tgz /var/npemwin/www
cd /var/npemwin/www
tar -xzf www.tgz
rm www.tgz
chown -R $user:$user *
cd htdocs/pub
#
# Note: If we reverse the test [ ! -e ...] && ....
# the statement will return 1 if rst exists and pkg_add will
# think that the script is returning an error.
#
[ -e "digatmos" ] || ln -s /var/npemwin/data/digatmos digatmos
[ -e "latest" ] || ln -s /var/npemwin/data/latest latest
[ -e "msg" ] || ln -s /var/npemwin/data/msg msg
[ -e "gempak" ] || ln -s /var/npemwin/data/gempak gempak
[ -e "track" ] || ln -s /var/npemwin/data/track track
[ -d "metarplot" ] || install -d -m 755 -o $user -g $user metarplot

#
# Display message
#
cat <<__EOF__
-
The configuration files are installed in /usr/local/etc/npemwin.

The installation process should have created various directories
under /var/npemwin, and already enabled the servers.conf file.
The program will run in this default installation without further
configuration.

Some common configuration options are indicated in various
sample configuration files in the /usr/local/etc/npemwin
directories and the defaults subdirectory.

The directory /usr/local/share/doc/npemwin contains
instructions for installing filters to process the files as they are
received.

An init script is installed to start automatically the program
at boot time, and stop when the system is shutdown.
It can also be started or stoped manually by executing

        /etc/init.d/npemwind start
        /etc/init.d/npemwind stop

The installation process should have created various directories
under /var/noaaport, and already enabled the rstfilter and dafilter.
In addition the web server is enabled by default, and can be
accesed as
	
	http://localhost:8016
or
	http://<hostname>:8016

The program will run in this default installation without further
configuration.

NOTE:
----
  As of version 2.0.r2, npemwind runs as user noaaport, and therefore
  the directories where the daemon writes must be owned by that user:

    chown -R noaaport:noaaport /var/npemwin
    chown -R noaaport:noaaport /var/log/npemwin
    chown -R noaaport:noaaport /var/run/npemwin

  During a fresh install, the package installer will create the
  directories with the appropriate ownership, but not during an upgrade.
  If you are uprading from an older version, those commands must be
  executed manually.
-
__EOF__
